# www/domotix_forever
### BEGIN INIT INFO
# Provides:             domotix
# Required-Start:
# Required-Stop:
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    Domotix plugin
### END INIT INFO

#export PATH=$PATH:/usr/local/bin
#export NODE_PATH=$NODE_PATH:/opt/node/lib/node_modules
export HOME=/home/pi/panStamp/domotix
export APP_FILE=$HOME/app.js

LOG_FILE=$HOME/logs/domotix_forever.log
LOG_OUT_FILE=$HOME/logs/domotix_out.log
LOG_ERR_FILE=$HOME/logs/domotix_err.log

# Required for socketstream
cd $HOME

case "$1" in
  start)
    rm $LOG_FILE
    rm $LOG_OUT_FILE
    rm $LOG_ERR_FILE
    forever start --uid "domotix" -l $LOG_FILE -o $LOG_OUT_FILE -e $LOG_ERR_FILE $APP_FILE
    ;;
  stop)
    forever stop domotix
    ;;
  restart)
    forever restart domotix
    ;;
  status)
    forever list | grep $APP_FILE | wc -l | sed -e 's/1/App Running/' | sed -e 's/0/App not started/'
    ;;
  *)

  echo "Usage: /etc/init.d/domotix {start|stop|restart|status}"
  exit 1
  ;;
esac
exit 0
